# supervisord-noproxy.conf
[program:station]
command=/scripts/station.bash

[supervisord]
nodaemon=true
